# src/variants/lpfrg/filesystem/root/home/pi/.config/autostart/chromium.desktop
[Desktop Entry]
Name=Chromium
Exec=chromium-browser http://localhost/ --kiosk --incognito --disable-pinch --no-first-run 
Type=application