# src/variants/lpfrg/filesystem/root/home/pi/.config/autostart/factoryreset.desktop
[Desktop Entry]
Type=Application
Name=Factory reset
Exec=/home/pi/scripts/factoryreset.py
Icon=
Comment=
