# src/variants/lpfrg/filesystem/root/home/pi/.config/autostart/touchrotate.desktop
[Desktop Entry]
Type=Application
Name=Rotate touch
Exec=/etc/dorotatetouch
Icon=
Comment=
